# systemd/certbot-renewal.timer
[Unit]
Description=Timer for Let's Encrypt certificate renewal (pet-info.link)
Requires=certbot-renewal.service

[Timer]
# Run daily at 3 AM
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
# Run 5 minutes after boot if the system was off during scheduled time
Persistent=true
# Randomize start time by up to 1 hour to avoid load spikes
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
