# systemd/pet-info.service.example
[Unit]
Description=Pet Info Web Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pet-info
Environment="ENV=prod"
Environment="APP_CERT_DIR=/opt/pet-info"

# Path to the compiled binary
ExecStart=/opt/pet-info/pet-info

# Automatic restart configuration
Restart=always
RestartSec=10

# Security hardening (optional - adjust as needed)
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pet-info

# Resource limits (optional)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
